# scripts/prismatron-network.service
[Unit]
Description=Prismatron Network Management Service
After=network.target NetworkManager.service
Wants=NetworkManager.service
BindsTo=NetworkManager.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'sleep 5; /usr/bin/python3 /opt/prismatron/scripts/network-startup.py'
ExecStop=/usr/bin/python3 /opt/prismatron/scripts/network-shutdown.py
User=root
Group=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
