# scripts/prismatron.service
[Unit]
Description=Prismatron LED Display System
Documentation=https://github.com/mwatson2/prismatron
After=network-online.target
Wants=network-online.target
# If using WLED, ensure network is truly ready
After=systemd-networkd-wait-online.service

[Service]
Type=forking
PIDFile=/run/prismatron/prismatron.pid
User=prismatron
Group=prismatron
WorkingDirectory=/mnt/dev/prismatron

# Environment setup
Environment="PYTHONPATH=/mnt/dev/prismatron"
Environment="PATH=/mnt/dev/prismatron/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="DISPLAY=:1"

# Pre-start cleanup (remove stale shared memory)
ExecStartPre=/bin/bash -c 'rm -f /dev/shm/prismatron_*'

# Main execution
ExecStart=/mnt/dev/prismatron/env/bin/python /mnt/dev/prismatron/main.py --daemon --config /mnt/dev/prismatron/config.json

# Graceful shutdown
ExecStop=/bin/kill -TERM $MAINPID
TimeoutStopSec=30
KillMode=mixed

# Restart policy
Restart=on-failure
RestartSec=10
StartLimitInterval=600
StartLimitBurst=3

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Allow writing to necessary directories
ReadWritePaths=/mnt/dev/prismatron/uploads
ReadWritePaths=/mnt/dev/prismatron/playlists
ReadWritePaths=/mnt/dev/prismatron/media
ReadWritePaths=/mnt/dev/prismatron/logs
ReadWritePaths=/mnt/dev/prismatron/thumbnails
ReadWritePaths=/run/prismatron
ReadWritePaths=/dev/shm
ReadWritePaths=/mnt/dev/prismatron

# Hardware access (for Jetson/LED control)
DeviceAllow=/dev/i2c-* rw
DeviceAllow=/dev/gpiochip* rw
DeviceAllow=/dev/video* rw
DeviceAllow=/dev/snd/* rw
SupplementaryGroups=video gpio i2c audio

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=8G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
